# _scripts/crontab
30 6 * * *     php /var/www/latestnews/latest/bin/console --env=prod news:fetchLatest
35 6 * * *     php /var/www/latestnews/latest/bin/console --env=prod news:fetchTitle
